# Build/supervisor/smartcambuddy.conf
[program:smartcambuddy]
command=/home/pi/.smartcambuddy_venv/bin/python /home/pi/SmartCamBuddy/smarcambuddy/read_sensors.py
autostart=true
autorestart=true
stdout_logfile=/home/pi/SmartCamBuddy/logs/main.log
redirect_stderr=true
